# setup/dashium.service
[Unit]
Description=Dashium

[Service]
WorkingDirectory=/dashium
ExecStart=sudo npm start
Restart=always

[Install]
WantedBy=multi-user.target
